# daemon/claude-usage-daemon.service
[Unit]
Description=Claude Usage Tracker Daemon (BLE)
Requires=bluetooth.target
After=bluetooth.target

[Service]
Type=simple
ExecStart=DAEMON_PATH
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
